# systemd/certbot.service
[Unit]
Description=Certbot
Documentation=https://certbot.eff.org/docs/
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Use writable directories for certbot configuration, work, and logs
# This is required in containerized/sandboxed environments where
# /etc/letsencrypt, /var/lib/letsencrypt, and /var/log/letsencrypt are read-only
ExecStart=/usr/bin/certbot renew --quiet \
    --config-dir /home/user/eas-station/certbot_data/config \
    --work-dir /home/user/eas-station/certbot_data/work \
    --logs-dir /home/user/eas-station/certbot_data/logs \
    --deploy-hook "systemctl reload nginx"
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
