# systemd/eas-station-audio.service
[Unit]
Description=EAS Station Audio Processing Service
Documentation=https://github.com/KR8MER/eas-station
After=network.target redis.service eas-station-sdr.service
Requires=redis.service
Wants=eas-station-sdr.service

[Service]
Type=simple
User=eas-station
Group=eas-station
# Additional group for audio device access
SupplementaryGroups=audio

WorkingDirectory=/opt/eas-station
Environment="PATH=/opt/eas-station/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/opt/eas-station"
EnvironmentFile=/opt/eas-station/.env

# Logging to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eas-station-audio

ExecStart=/opt/eas-station/venv/bin/python eas_monitoring_service.py

# Service management
Restart=always
RestartSec=10s
KillMode=mixed
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/eas-station /var/log/eas-station

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
