# systemd/eas-station-poller.service
[Unit]
Description=EAS Station Alert Poller (Unified - NOAA + IPAWS)
Documentation=https://github.com/KR8MER/eas-station
PartOf=eas-station.target
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=eas-station
Group=eas-station
WorkingDirectory=/opt/eas-station
Environment="PATH=/opt/eas-station/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/opt/eas-station"
Environment="PYTHONUNBUFFERED=1"
# Pi 5 requires lgpio for GPIO access - prevents gpiozero auto-detection hangs
Environment="GPIOZERO_PIN_FACTORY=lgpio"
# NOTE: CAP_POLLER_MODE is NOT set - unified poller polls all sources (NOAA + IPAWS)
# Configure sources via .env: IPAWS_CAP_FEED_URLS, CAP_ENDPOINTS, NOAA_USER_AGENT
# Application defaults to /opt/eas-station/.env (can override with CONFIG_PATH env var)
EnvironmentFile=/opt/eas-station/.env
ExecStart=/opt/eas-station/venv/bin/python /opt/eas-station/poller/cap_poller.py --continuous
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eas-poller

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target eas-station.target
