# systemd/eas-station-sdr.service
[Unit]
Description=EAS Station SDR Hardware Service
Documentation=https://github.com/KR8MER/eas-station
PartOf=eas-station.target
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=eas-station
Group=eas-station
# Additional groups for USB/hardware access
SupplementaryGroups=plugdev dialout

WorkingDirectory=/opt/eas-station
Environment="PATH=/opt/eas-station/venv-sdr/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# venv-sdr is created with --system-site-packages so it can access python3-soapysdr from apt
# This eliminates the need for complex PYTHONPATH hacks - SoapySDR is directly importable
# The web service uses a separate isolated venv to avoid gunicorn/gevent conflicts
Environment="PYTHONPATH=/opt/eas-station"
# SoapySDR plugin path for hardware driver modules (architecture-specific)
Environment="SOAPY_SDR_PLUGIN_PATH=/usr/lib/aarch64-linux-gnu/SoapySDR/modules0.8:/usr/lib/arm-linux-gnueabihf/SoapySDR/modules0.8:/usr/lib/x86_64-linux-gnu/SoapySDR/modules0.8:/usr/lib/SoapySDR/modules0.8"
# Suppress ALSA/PulseAudio/Jack errors from soapysdr-module-audio (if installed)
# soapysdr-module-audio tries to probe ALSA/Pulse/Jack devices even though we only use SDR hardware
# Setting these environment variables suppresses the error messages without affecting SDR functionality
Environment="ALSA_CONFIG_PATH=/dev/null"
Environment="PULSE_SERVER=/dev/null"
Environment="SOAPY_SDR_LOG_LEVEL=WARNING"
# Pi 5 requires lgpio for GPIO access - prevents gpiozero auto-detection hangs
Environment="GPIOZERO_PIN_FACTORY=lgpio"
EnvironmentFile=/opt/eas-station/.env

# Logging to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eas-station-sdr

ExecStart=/opt/eas-station/venv-sdr/bin/python sdr_hardware_service.py

# Service management
Restart=always
RestartSec=10s
KillMode=mixed
TimeoutStopSec=30

# Security settings - relaxed for USB/SDR hardware access
# SDR hardware requires access to /dev/bus/usb and /sys for device enumeration
NoNewPrivileges=true
PrivateTmp=true

# Resource limits - USB streaming requires locked memory for DMA
LimitNOFILE=65536
LimitMEMLOCK=infinity
MemoryMax=256M

[Install]
WantedBy=multi-user.target eas-station.target
